# config/pdf-celery-beat.service
[Unit]
Description=PDF Tools Celery Beat Scheduler
After=network.target redis.service

[Service]
Type=simple
User=pdfapi
Group=pdfapi
WorkingDirectory=/home/pdfapi/pdf-tools-backend
Environment="PATH=/home/pdfapi/pdf-tools-backend/venv/bin"
EnvironmentFile=/home/pdfapi/pdf-tools-backend/.env
ExecStart=/home/pdfapi/pdf-tools-backend/venv/bin/celery -A celery_app beat \
    --loglevel=info \
    --logfile=/var/log/pdf-api/celery-beat.log \
    --pidfile=/var/run/celery/celery-beat.pid

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
