# config/pdf-celery.service
[Unit]
Description=PDF Tools Celery Worker
After=network.target redis.service

[Service]
Type=forking
User=pdfapi
Group=pdfapi
WorkingDirectory=/home/pdfapi/pdf-tools-backend
Environment="PATH=/home/pdfapi/pdf-tools-backend/venv/bin"
EnvironmentFile=/home/pdfapi/pdf-tools-backend/.env
ExecStart=/home/pdfapi/pdf-tools-backend/venv/bin/celery -A celery_app worker \
    --loglevel=info \
    --logfile=/var/log/pdf-api/celery.log \
    --pidfile=/var/run/celery/celery.pid \
    --concurrency=4 \
    --max-tasks-per-child=100

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
